# config/couchpotato.ini
[program:CouchPotatoServer]
command=python /opt/couchpotato/CouchPotato.py --data_dir /config
user=root