# config/cron.ini
[program:CronJob]
command=/usr/sbin/crond -f -d 8
user=root